# SQA/supervisord.conf
[supervisord]
nodaemon=true

[program:docker-compose]
command=/app/script_dockerfile.sh
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
autorestart=true
